# nad_receive.service
[Unit]
Description=NAD Receiver Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nad
Group=nad
SupplementaryGroups=dialout
WorkingDirectory=/opt/nad
ExecStart=/opt/nad/venv/bin/python /opt/nad/receiver.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
